# supervisord.conf
[program:php]
directory=/w/php/furion/public
command=php -S 0.0.0.0:80
